# build/scripts/sealion
#! /bin/sh

### BEGIN INIT INFO
# Provides:          sealion
# Required-Start:       
# Required-Stop:     
# Should-Start:      
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: SeaLion agent
# Description:       SeaLion is cloud based server debugging tool.
#                    This is SeaLion's Agent to send monitoring data to SeaLion server
# 
### END INIT INFO

BASEDIR=$(dirname $0)
BASEDIR=${BASEDIR%/}
python $BASEDIR/../sealion.py $1

