# build/res/scripts/sealion
#!/usr/bin/env bash

### BEGIN INIT INFO
# Provides: sealion
# Required-Start:       
# Required-Stop:     
# Should-Start:      
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: SeaLion Agent
# Description: SeaLion is cloud based server debugging tool. This is SeaLion's agent to send monitoring data to SeaLion server
### END INIT INFO

BASEDIR=$(readlink -f "$0")
BASEDIR=$(dirname "$BASEDIR")
BASEDIR=${BASEDIR%/}
python "$BASEDIR/../../src/sealion.py" $1
exit $?

